# deploy/systemd/anonimus.service
[Unit]
Description=anonimus_chat production Docker stack
Documentation=https://github.com/flykby/anonimus_chat/blob/main/docs/deploy.md
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/anonimus_chat
EnvironmentFile=-/opt/anonimus_chat/.env
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml down
ExecReload=/opt/anonimus_chat/scripts/deploy.sh
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
